# ansible-common/files/common/profile.d/awscli.sh
if [ -r /etc/awscli.conf ]; then
    export AWS_CONFIG_FILE=/etc/awscli.conf
fi
